# image-metadata/init/supervisord.ini
[supervisord]
nodaemon=true

[program:hdfsmount]
command=/root/set_mount.sh
startsecs=0
priority=100
redirect_stderr=true
stdout_logfile=/var/log/mount_initialization.log
autorestart=false

[program:sshd]
command=/usr/sbin/sshd -D
