# docker-base/init/supervisord.ini
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:dockerd]
command=/usr/bin/dockerd -D --experimental -H unix:///var/run/docker.sock -H tcp://0.0.0.0:2375
